# nginx/nginx+nodejs/supervisord.d/supervisord.nginx.conf
[program:nginx]
command = /usr/sbin/nginx -g "daemon off;"
autostart = true
autorestart = true

stdout_logfile = none
stderr_logfile = /var/log/nginx.err

stderr_events_enabled = true
