# nginx/nginx+nodejs/supervisord.d/supervisord.nodejs.conf
[program:nodejs]
command = npm start
autostart = true
autorestart = true
directory = /app

stdout_logfile = /var/log/nodejs.log
stderr_logfile = /var/log/nodejs.err
